# systemd/updateDNS.service
[Unit]
Description=Update DDNS IP at startup using ddclient

[Service]
ExecStart=/usr/sbin/ddclient --debug

[Install]
WantedBy=multi-user.target